# cloudkey_stock/sbin/ubnt-systemhub.service
[Unit]
Description=UBNT SystemHub

[Service]
Type=notify
Restart=always
ExecStart=/usr/sbin/ubnt-systemhub

[Install]
WantedBy=multi-user.target
